# gunicorn.conf
description "Gunicorn application server handling aerabi.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid root
setgid www-data
chdir /root/aerabi/aerabi

exec venv/bin/gunicorn --workers 3 --bind unix:/root/aerabi/gunicorn.sock aerabi.wsgi:application
